# infra/pi-image/systemd/showroom-agent.service
[Unit]
Description=Showroom device agent
After=network-online.target
Wants=network-online.target

[Service]
EnvironmentFile=/etc/showroom-agent/config.env
ExecStart=/usr/local/bin/showroom-agent
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target

